# home/config/systemd/user/conky.service
[Unit]
Description=Conky System Monitor
After=graphical.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'sleep 5 && /usr/bin/conky --quiet --config=$HOME/.conkyrc'
Restart=on-failure
RestartSec=10
Environment=DISPLAY=:0

[Install]
WantedBy=default.target

